# tree/gen/home/etc/service/chanal/run
#!/bin/sh
exec /usr/local/bin/py2 ~/p/one-off/pgnanal/py/fics.py
