# tree/gen/etc/service/memoplex/run
#!/bin/sh
exec sudo -u danl memoplex -s
