# tree/gen/etc/service/memoplex/run
#!/bin/sh
exec sudo -u danl /usr/local/bin/memoplex -s -p dan2 > /var/log/memoplex 2>&1
