# tree/gen/home/service/bunny1/run
#!/bin/sh
cd ~/i/mod/bunny1/wtf && exec python a.py
