# old-tree/gen/service/memoplex/run
#!/bin/sh
exec sudo -u danl /usr/local/bin/memoplex -s > /var/log/memoplex 2>&1
