# tree/gen/etc/service/pulseaudio/run
#!/bin/sh
exec pulseaudio --realtime --system
