# tree/gen/etc/service/rec-midi/run
#!/bin/sh
exec sudo -H -u danl /etc/service/rec_midi/run.danl
