# tree/gen/home/etc/service/skim-zim/run
#!/bin/sh
exec ~/p/one-off/skim-zim/serve ~/data/enwikt.zim
