# tree/gen/home/etc/service/skim-zim/run
#!/bin/sh
exec ~/p/o/skim-zim/serve ~/data/wikt/en.zim ~/data/wikt/es.zim \
  ~/data/wiki/en.zim
